# files/cleverbreeze-ui.service
[Unit]
Description=Clever Breeze UI Service

[Service]
ExecStart= /usr/bin/env uwsgi --ini /root/cleverbreeze/system/uwsgi.ini
Restart=always
KillSignal=SIGTERM
Type=notify
NotifyAccess=all
User=root
Group=root

[Install]
WantedBy=multi-user.target
